# services/mesh-guardian-backup.service
[Unit]
Description=AILinux Mesh Guardian - Self-Healing System
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=backupuser
Group=backupuser
WorkingDirectory=/home/backupuser/triforce
ExecStart=/home/backupuser/triforce/.venv/bin/python /home/backupuser/triforce/scripts/mesh-guardian.py --interval 30
Restart=always
RestartSec=10
StandardOutput=append:/var/log/mesh-guardian.log
StandardError=append:/var/log/mesh-guardian.log

[Install]
WantedBy=multi-user.target
