# systemd/federation-node.service
[Unit]
Description=AILinux Federation Node Agent
After=network.target ollama.service

[Service]
Type=simple
User=zombie
WorkingDirectory=/home/zombie/triforce
Environment="PATH=/home/zombie/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /home/zombie/triforce/federation_node.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
